# duckdns-updater.service
[Unit]
Description=DuckDNS IP Updater for NOCTIS PRO
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/workspace
ExecStart=/workspace/update_duckdns.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target