# noctispro-ngrok-current.service
[Unit]
Description=NoctisPro Ngrok Tunnel (Current Setup)
After=network-online.target noctispro-production-current.service
Wants=network-online.target
Requires=noctispro-production-current.service

[Service]
Type=simple
Restart=always
RestartSec=10
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace
Environment=HOME=/home/ubuntu
ExecStartPre=/bin/sleep 15
ExecStart=/usr/bin/ngrok http 8000 --log=stdout
ExecStop=/usr/bin/pkill -f ngrok
StandardOutput=append:/workspace/logs/ngrok.log
StandardError=append:/workspace/logs/ngrok-error.log
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target