# noctispro-production-complete.service
[Unit]
Description=NoctisPro Complete Production System with Full DICOM Viewer
Documentation=https://noctispro.com/docs
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target postgresql.service redis-server.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/workspace
Environment=PATH=/workspace/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/var/lib/www-data
EnvironmentFile=-/workspace/.env.production
EnvironmentFile=-/workspace/.env

# Ensure production DICOM viewer is used
Environment=DJANGO_SETTINGS_MODULE=noctis_pro.settings_production
Environment=DEBUG=False
Environment=USE_SQLITE=false
Environment=DISABLE_REDIS=false
Environment=USE_DUMMY_CACHE=false

# Pre-start health checks
ExecStartPre=/bin/bash -c 'systemctl is-active postgresql || systemctl start postgresql'
ExecStartPre=/bin/bash -c 'systemctl is-active redis-server || systemctl start redis-server'
ExecStartPre=/bin/sleep 10

# Database and static files preparation
ExecStartPre=/bin/bash -c 'cd /workspace && source venv/bin/activate && python manage.py migrate --noinput --settings=noctis_pro.settings_production'
ExecStartPre=/bin/bash -c 'cd /workspace && source venv/bin/activate && python manage.py collectstatic --noinput --settings=noctis_pro.settings_production'

# Create necessary directories
ExecStartPre=/bin/bash -c 'mkdir -p /workspace/media/dicom /workspace/staticfiles /workspace/logs'
ExecStartPre=/bin/bash -c 'chown -R www-data:www-data /workspace/media /workspace/staticfiles'

# Start the main application with production DICOM viewer
ExecStart=/workspace/venv/bin/daphne -b 0.0.0.0 -p 8000 noctis_pro.asgi:application

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30
KillMode=mixed

# Restart configuration
Restart=always
RestartSec=15
TimeoutStartSec=300

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noctispro-production

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/workspace/media /workspace/staticfiles /workspace/logs /tmp
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Health monitoring
WatchdogSec=60

[Install]
WantedBy=multi-user.target