# noctispro-professional.service
[Unit]
Description=NoctisPro Professional Medical Imaging System
Documentation=https://noctispro.com/docs
After=network.target
Wants=network.target

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace/noctis_pro_deployment
Environment=DJANGO_SETTINGS_MODULE=noctis_pro.settings_production
Environment=PYTHONPATH=/workspace/noctis_pro_deployment
ExecStartPre=/bin/bash -c 'cd /workspace/noctis_pro_deployment && source venv/bin/activate && python manage.py collectstatic --noinput --clear'
ExecStartPre=/bin/bash -c 'cd /workspace/noctis_pro_deployment && source venv/bin/activate && python manage.py migrate'
ExecStart=/bin/bash -c 'cd /workspace/noctis_pro_deployment && source venv/bin/activate && gunicorn noctis_pro.wsgi:application --bind 0.0.0.0:8000 --workers 4 --timeout 300 --max-requests 1000 --preload'
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noctispro-professional

# Professional resource limits for medical imaging
LimitNOFILE=65536
LimitNPROC=4096

# Professional security settings
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target