# noctispro-init
#!/bin/bash
### BEGIN INIT INFO
# Provides:          noctispro
# Required-Start:    $network $remote_fs $syslog
# Required-Stop:     $network $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: NoctisPro Medical Imaging System
# Description:       Starts NoctisPro Django application with ngrok tunnel
### END INIT INFO

# Change to the workspace directory
cd /workspace

case "$1" in
    start)
        echo "Starting NoctisPro service..."
        /workspace/noctispro_service.sh start
        ;;
    stop)
        echo "Stopping NoctisPro service..."
        /workspace/noctispro_service.sh stop
        ;;
    restart)
        echo "Restarting NoctisPro service..."
        /workspace/noctispro_service.sh restart
        ;;
    status)
        /workspace/noctispro_service.sh status
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0