# noctispro-ngrok.service
[Unit]
Description=NoctisPro Ngrok Tunnel
After=network.target noctispro.service
Wants=noctispro.service

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace
Environment="HOME=/home/ubuntu"
EnvironmentFile=-/workspace/.env.ngrok
ExecStartPre=/bin/bash -c 'source /workspace/.env.ngrok 2>/dev/null || true'
ExecStart=/bin/bash -c 'cd /workspace && source .env.ngrok 2>/dev/null || true && if [ "${NGROK_USE_STATIC:-false}" = "true" ]; then if [ ! -z "${NGROK_DOMAIN:-}" ]; then /usr/local/bin/ngrok start noctispro-domain --log=stdout; elif [ ! -z "${NGROK_SUBDOMAIN:-}" ]; then /usr/local/bin/ngrok start noctispro-static --log=stdout; else /usr/local/bin/ngrok start ${NGROK_TUNNEL_NAME:-noctispro-http} --log=stdout; fi; else /usr/local/bin/ngrok start ${NGROK_TUNNEL_NAME:-noctispro-http} --log=stdout; fi'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target