# noctispro-ngrok.service
[Unit]
Description=NoctisPro Ngrok Tunnel
After=network.target noctispro-django.service
Requires=noctispro-django.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace
ExecStart=/workspace/ngrok http 8000 --hostname=colt-charmed-lark.ngrok-free.app --log=stdout
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target