# backup_misc/noctispro-reliable.service
[Unit]
Description=NoctisPro DICOM PACS System - Reliable Boot Service
Documentation=https://github.com/noctispro/pacs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/workspace
Environment=PYTHONPATH=/workspace
Environment=DJANGO_SETTINGS_MODULE=noctis_pro.settings
Environment=PYTHONUNBUFFERED=1

# Pre-start checks and setup
ExecStartPre=/bin/bash -c 'cd /workspace && python3 manage.py collectstatic --noinput --clear || true'
ExecStartPre=/bin/bash -c 'cd /workspace && python3 manage.py migrate --noinput || true'

# Main service command
ExecStart=/bin/bash -c 'cd /workspace && python3 -m gunicorn noctis_pro.wsgi:application --bind 0.0.0.0:8000 --workers 3 --timeout 120 --access-logfile /workspace/gunicorn_access.log --error-logfile /workspace/gunicorn_error.log --daemon && /workspace/start_ngrok.sh'

# Health check and restart policy
Restart=always
RestartSec=10
StartLimitBurst=5

# Resource limits
MemoryLimit=2G
CPUQuota=200%

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/workspace /tmp /var/tmp
ProtectHome=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noctispro

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target