# ops/duckdns-update.service
[Unit]
Description=DuckDNS IP Updater for Noctis
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/noctis/duckdns.env
ExecStart=/usr/bin/bash -lc '/usr/local/bin/duckdns-update.sh'

[Install]
WantedBy=multi-user.target