# ops/noctis-celery.service
[Unit]
Description=Noctis Pro Celery Worker
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/noctis/noctis.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/bash -lc '"$VENV_DIR/bin/celery" -A "$CELERY_APP" worker -l "$CELERY_LOGLEVEL"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target