# ops/noctis-web.service
[Unit]
Description=Noctis Pro ASGI (Daphne)
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/noctis/noctis.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/bash -lc '"$VENV_DIR/bin/python" "$APP_DIR/manage.py" migrate --noinput'
ExecStartPre=/usr/bin/bash -lc '"$VENV_DIR/bin/python" "$APP_DIR/manage.py" collectstatic --noinput'
ExecStart=/usr/bin/bash -lc '"$VENV_DIR/bin/daphne" -b "$HOST" -p "$PORT" "$ASGI_APP"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target