# ops/noctis-web.service
[Unit]
Description=Noctis Pro ASGI (Daphne)
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/workspace
Environment=DJANGO_SETTINGS_MODULE=noctis_pro.settings
Environment=REDIS_URL=redis://127.0.0.1:6379/0
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/bash -lc '/workspace/venv/bin/python /workspace/manage.py migrate --noinput'
ExecStartPre=/usr/bin/bash -lc '/workspace/venv/bin/python /workspace/manage.py collectstatic --noinput'
ExecStart=/workspace/venv/bin/daphne -b 127.0.0.1 -p 8000 noctis_pro.asgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target