# ops/noctis-dicom.service
[Unit]
Description=Noctis Pro DICOM Receiver (SCP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/workspace
Environment=PYTHONUNBUFFERED=1
ExecStart=/workspace/venv/bin/python /workspace/dicom_receiver.py --port 11112 --aet NOCTIS_SCP
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target