# ops/noctis-dicom.service
[Unit]
Description=Noctis Pro DICOM Receiver (SCP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/noctis/noctis.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/bash -lc '"$VENV_DIR/bin/python" "$APP_DIR/dicom_receiver.py" --port "$DICOM_PORT" --aet "$DICOM_AET"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target